# etc/systemd/system/weekly-maintenance.service
[Unit]
Description=Weekly system maintenance (updates + repair)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/weekly-maintenance.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
